# templates/default/upstart.conf.erb
start on runlevel [2345]
stop on runlevel [016]

respawn

env PORT=5001

setuid <%= @user %>

chdir <%= @path %>

limit nofile 500000 500000

exec env BUNDLE_PATH=<%= @bundle_path %> bash -lx bundle exec ruby server.rb --vv -c <%= @config_file %>
